# roles/git-nagcollect/templates/nagcollect.cron.j2
# /etc/cron.d/nagcollect
# ---
#  {{ ansible_managed }}
# ---

*/15 * * * *	root	/usr/local/sbin/nagcollect >/dev/null
